# com.user.ai-news-summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.ai-news-summary</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/user/ai-news-summary/venv/bin/python</string>
        <string>/Users/user/ai-news-summary/src/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/user/ai-news-summary</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>GEMINI_API_KEY</key>
        <string>REDACTED_GEMINI_API_KEY</string>
        <key>YOUTUBE_API_KEY</key>
        <string>REDACTED_YOUTUBE_API_KEY</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/user/ai-news-summary/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/user/ai-news-summary/launchd.error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
